# machines/proxy-server/services/tor.service.d/allow.conf
[Service]
ReadWriteDirectories=/run/socket/tor
