# machines/samba/service/flush-config.service
[Service]
EnvironmentFile=/etc/environment
ExecStartPre=/bin/bash /opt/scripts/prestart.sh
