# machines/samba/service/flush-config.service
[Service]
ExecStartPre=/bin/bash /opt/scripts/prestart.sh
