# machines/webservice/service/nextcloud-cron.timer
[Unit]
Description=nextcloud cron.php
Wants=php-fpm.service

[Timer]
OnBootSec=10min
OnUnitInactiveSec=1h

[Install]
WantedBy=multi-user.target

